# trellm.service
[Unit]
Description=TreLLM - Trello to Claude Code automation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=TRELLM_USER
ExecStart=TRELLM_VENV/bin/trellm
Restart=on-failure
RestartSec=10
WorkingDirectory=TRELLM_DIR

# Ensure PATH includes common locations for claude, node, etc.
Environment=PATH=TRELLM_HOME/.local/bin:TRELLM_HOME/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=TRELLM_HOME

[Install]
WantedBy=multi-user.target
